# docker/app/supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/var/run/supervisord.pid

[program:php-fpm]
command=php-fpm --nodaemonize
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:queue]
directory=/var/www/html
command=sh -lc 'while [ ! -f artisan ]; do echo "[queue] waiting for app..."; sleep 2; done; php artisan queue:work --sleep=1 --tries=3 --max-time=3600'
user=www-data
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:schedule]
directory=/var/www/html
command=sh -lc 'while [ ! -f artisan ]; do echo "[schedule] waiting for app..."; sleep 2; done; php artisan schedule:work'
user=www-data
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

